# project-rpi/logger/daemonlogger.service
[unit]
Description=app
[Service]
ExecStart=/home/usef/project/logger/build/app
[Install]
WantedBy=multi-user.target
